# crontab
# ─── Access100 API Cron Jobs ─────────────────────────────────────────
# All times in UTC. Hawaii (HST) is UTC-10.
# Each job sources /etc/cron-env.sh (written by docker-entrypoint.sh)
# so that Docker environment variables are available to PHP.
# Jobs run as www-data (same user as Apache workers) — never root.

SHELL=/bin/bash

# Scrape: eHawaii RSS feed polling (hourly)
0 * * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/scrape.php >> /var/log/access100-scrape.log 2>&1

# Scrape: NCO neighborhood board meetings (daily, 6 AM HST = 4 PM UTC)
0 16 * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/scrape-nco.php >> /var/log/access100-scrape-nco.log 2>&1

# Scrape: Honolulu boards & commissions (daily, 6:02 AM HST = 4:02 PM UTC)
2 16 * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/scrape-honolulu-boards.php >> /var/log/access100-scrape-hnl-boards.log 2>&1

# Scrape: Maui County Legistar (hourly, 6 min after eHawaii)
6 * * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/scrape-maui-legistar.php >> /var/log/access100-scrape-maui.log 2>&1

# Notify: change detection & immediate notifications (hourly, 5 min after scrape)
5 * * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/notify.php >> /var/log/access100-notify.log 2>&1

# Daily digest (7 AM HST = 5 PM UTC)
0 17 * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/digest.php >> /var/log/access100-digest.log 2>&1

# Weekly digest (Monday 7 AM HST = Monday 5 PM UTC)
0 17 * * 1 www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/weekly-digest.php >> /var/log/access100-weekly.log 2>&1

# Meeting reminders (7 AM HST = 5 PM UTC)
0 17 * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/reminder.php >> /var/log/access100-reminder.log 2>&1

# AI topic classification (2 AM HST = 12 PM UTC)
0 12 * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/classify-topics.php >> /var/log/access100-classify.log 2>&1

# Cleanup: rate limit files + old queue entries (3:30 AM HST = 1:30 PM UTC)
30 13 * * * www-data . /etc/cron-env.sh && /usr/local/bin/php /var/www/html/api/cron/cleanup.php >> /var/log/access100-cleanup.log 2>&1

# Required trailing newline for cron
